# systemd/awesomesauce-digest.timer
[Unit]
Description=Timer for Hackernews digest scraper and posting tool

[Timer]
OnCalendar=*-*-* 14:30:00

[Install]
WantedBy=timers.target
